# systemd/system/elton.timer
[Unit]
Description=Periodically Run Elton

[Timer]
OnUnitActiveSec=2d
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
